# deploy/ffxiv-flask.service
[Unit]
Description=FF14 Party Planner — Flask API
After=network.target

[Service]
Type=simple
User=YOUR_USER
Group=YOUR_USER
WorkingDirectory=/path/to/ffxivpp
ExecStart=/path/to/ffxivpp/.venv/bin/gunicorn wsgi:app --bind 127.0.0.1:7111 --workers 2 --access-logfile /tmp/flask-access.log
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
